# agent_packages/resolver/debian/deb-package/DEBIAN/prerm
#!/bin/sh
systemctl stop cxoneflow-resolver-agent
systemctl disable cxoneflow-resolver-agent
