# agent_packages/resolver/debian/deb-package/DEBIAN/prerm
#!/bin/sh
echo Stopping the 'cxoneflow-resolver-agent' service.
systemctl stop cxoneflow-resolver-agent
systemctl disable cxoneflow-resolver-agent
